# scripts/installDrivers.service
[Unit]
Description=install all drivers 

[Service]
ExecStart=/home/ubuntu/installDrivers.sh

[Install]
WantedBy=default.target